# crontab
*/30 * * * * /usr/bin/getmail --rcfile getmail.gmail
